# parking.conf
# put this file to /etc/init/
start on started apache2
exec /etc/init.d/parking start
