# queue.service
[Unit]
Description=Parking Queue Service
Requires=network.target

[Service]
User=parking
WorkingDirectory=/home/parking/parking-system
ExecStart=/usr/bin/php artisan queue:work
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
